# files/etc/rc.local

# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.
echo "0x19d2 0x0532" > /sys/bus/usb-serial/drivers/generic/new_id
sleep 1
gcom -d /dev/ttyUSB0 -s /proc/self/fd/0 << EOF
opengt
 set com 115200n81
 set comecho off
 set senddelay 0.02
 waitquiet 0.2 0.2
 flash 0.1

exit 0
